# deploy/upstart/vagrant.conf
#!upstart
description "boiler api server"
author      "zach <zach@zed.io>"

start on started networking
stop on shutdown

respawn
respawn limit 5 5

env PROG_NAME=api
env INSTALL_DIR=/vagrant
env LOG_PATH=/vagrant/log/api.log

setgid vagrant
setuid vagrant

script
    echo $$ > $INSTALL_DIR/log/$PROG_NAME.pid
    node --harmony --stack-trace-limit=1000 ${INSTALL_DIR}/bin/$PROG_NAME >> ${LOG_PATH} 2>&1
end script
